# battery-notifications-h3f/battery-notifications.service
[Unit]
Description=Battery notifications

[Service]
Type=simple
ExecStart=/usr/bin/battery-notifications
